# deploy/systemd/pattern-scanner.service
[Unit]
Description=Pattern Scanner — Binance USDT永续合约K线形态识别引擎
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/projects/crypto
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/ubuntu/projects/crypto
EnvironmentFile=-/home/ubuntu/projects/crypto/.env
ExecStart=/home/ubuntu/projects/crypto/.venv/bin/python -m pattern_scanner.scheduler
Restart=always
RestartSec=10
# Web dashboard runs on port 8082 — ensure firewall allows it

[Install]
WantedBy=multi-user.target
